# config/services/mupi_splash.service
[Unit]
Description=Simple Boot Splash

[Service]
ExecStart=/usr/bin/fbv /boot/splash.png

[Install]
WantedBy=basic.target