# config/services/mupi_hat.service
[Unit]
Description=MuPiHAT Service
Before=basic.target
After=local-fs.target sysinit.target
DefaultDependencies=no

[Service]
Type=simple
RemainAfterExit=true
ExecStart=/usr/local/bin/mupibox/mupihat_automation.sh
ExecStart=/usr/bin/python3 -B /usr/local/bin/mupibox/mupihat.py -j /tmp/mupihat.json
ExecStop=pkill -f -9 mupihat_automation.sh
ExecStop=pkill -f -9 mupihat.py
KillMode=process
PIDFile=/run/mupi_hat.pid

[Install]
WantedBy=basic.target