# config/services/mupi_splash.service
[Unit]
Description=Simple Boot Splash

[Service]
ExecStart=/usr/local/bin/mupibox/./splash_screen.sh

[Install]
WantedBy=basic.target