# config/services/mupi_idle_shutdown.service
[Unit]
Description=Idle Shutdown Timer

[Service]
ExecStart=/usr/local/bin/mupibox/./idle_shutdown.sh

[Install]
WantedBy=basic.target