# config/services/mupi_wifi.service
[Unit]
Description=Add Wifi from Display Dialog

[Service]
ExecStart=/usr/local/bin/mupibox/./add_wifi.sh

[Install]
WantedBy=basic.target