# config/services/mupi_powerled.service
[Unit]
Description=Power LED Control

[Service]
Type=simple
RemainAfterExit=true
ExecStart=/usr/local/bin/mupibox/./mupi_start_led.sh
KillMode=process

[Install]
WantedBy=basic.target