# config/services/mupi_hat.service
[Unit]
Description=MuPiHAT Service
Before=basic.target
After=local-fs.target sysinit.target
DefaultDependencies=no

[Service]
Type=simple
WorkingDirectory=/usr/local/bin/mupibox/
User=root
ExecStart=/usr/bin/python3 -B /usr/local/bin/mupibox/mupihat.py -j /tmp/mupihat.json
Restart=on-failure

[Install]
WantedBy=basic.target