# config/services/mupi_check_internet.service
[Unit]
Description=Checks Internet Connection

[Service]
User=dietpi
Group=dietpi
ExecStart=/usr/local/bin/mupibox/./check_network.sh

[Install]
WantedBy=multi-user.target