# config/services/mupi_powerled.service
[Unit]
Description=Power LED Control

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/bin/mupibox/./mupi_start_led.sh
ExecStop=/usr/local/bin/mupibox/./mupi_stop_led.sh

[Install]
WantedBy=basic.target