# config/services/mupi_check_monitor.service
[Unit]
Description=Checks Blank Monitor

[Service]
ExecStart=/usr/local/bin/mupibox/./get_monitor.sh

[Install]
WantedBy=default.target